# xmas.service
[Unit]
Description=PiHut 3D Xmas Tree Driver

[Service]
ExecStart=/usr/local/bin/xmas.py

[Install]
WantedBy=multi-user.target

